# goldeaf-backup.service
[Unit]
Description=GOLDEAF Bank Backup Service
After=network.target

[Service]
Type=oneshot
ExecStart=/var/www/goldeaf_bank/backup.sh
User=www-data
Group=www-data
